# .config/crontab
# cleanup download folder
*/10 * * * * $HOME/.local/bin/obliviate $HOME/Downloads
