# .config/crontab
# cleanup download folder
0 10 * * * $HOME/.local/bin/obliviate $HOME/Downloads
